# deployment/supervisord-allinone.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid
user=root
environment=OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s"

# ============================================================================
# Supervisord Web UI - Access at http://localhost:9002
# ============================================================================
[inet_http_server]
port=*:9002
username=admin
password=admin

[supervisorctl]
serverurl=http://localhost:9002
username=admin
password=admin

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ============================================================================
# Infrastructure Services
# ============================================================================

[program:mongodb]
command=/usr/local/bin/mongod --dbpath=/data/mongodb --bind_ip_all --quiet --logpath=/dev/null
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nats]
command=/usr/local/bin/nats-server -js --store_dir=/data/nats --addr=0.0.0.0 --http_port=8222
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:minio]
command=/usr/local/bin/minio server /data/minio --address=:9000 --console-address=:9001
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=MINIO_ROOT_USER="minioadmin",MINIO_ROOT_PASSWORD="minioadmin"
priority=30

# ============================================================================
# Management UIs
# ============================================================================

[program:mongo-express]
command=/usr/local/bin/mongo-express
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=ME_CONFIG_MONGODB_URL="mongodb://localhost:27017/",ME_CONFIG_MONGODB_ENABLE_ADMIN="true",ME_CONFIG_BASICAUTH_USERNAME="admin",ME_CONFIG_BASICAUTH_PASSWORD="admin",ME_CONFIG_SITE_BASEURL="/",VCAP_APP_HOST="0.0.0.0",VCAP_APP_PORT="8081"
priority=100

# ============================================================================
# Application Services
# ============================================================================

[program:api]
command=python -m uvicorn src.api.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",S3_ENDPOINT_URL="http://localhost:9000",S3_ACCESS_KEY_ID="minioadmin",S3_SECRET_ACCESS_KEY="minioadmin",S3_BUCKET_NAME="creative-assets",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",API_ENABLE_READINESS_PROBE="false",LOG_LEVEL="INFO"
priority=200

[program:web]
command=streamlit run src/web/app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=API_BASE_URL="http://localhost:8000",LOG_LEVEL="INFO"
priority=210

[program:context-enricher]
command=python -m src.context_enricher.main
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",CONTEXT_ENRICHER_ENABLE_READINESS_PROBE="false",LOG_LEVEL="INFO"
priority=300

[program:creative-generator]
command=python -m src.creative_generator.main
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",CREATIVE_GENERATOR_ENABLE_READINESS_PROBE="false",LOG_LEVEL="INFO"
priority=310

[program:image-generator]
command=python -m src.image_generator.main
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",S3_ENDPOINT_URL="http://localhost:9000",S3_ACCESS_KEY_ID="minioadmin",S3_SECRET_ACCESS_KEY="minioadmin",S3_BUCKET_NAME="creative-assets",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",IMAGE_GENERATOR_ENABLE_READINESS_PROBE="false",LOG_LEVEL="INFO"
priority=320

[program:brand-composer]
command=python -m src.brand_composer.main
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",S3_ENDPOINT_URL="http://localhost:9000",S3_EXTERNAL_ENDPOINT_URL="http://localhost:9000",S3_ACCESS_KEY_ID="minioadmin",S3_SECRET_ACCESS_KEY="minioadmin",S3_BUCKET_NAME="creative-assets",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",BRAND_COMPOSER_ENABLE_READINESS_PROBE="false",LOG_LEVEL="INFO"
priority=330

[program:text-overlay]
command=python -m src.text_overlay.main
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",S3_ENDPOINT_URL="http://localhost:9000",S3_EXTERNAL_ENDPOINT_URL="http://localhost:9000",S3_ACCESS_KEY_ID="minioadmin",S3_SECRET_ACCESS_KEY="minioadmin",S3_BUCKET_NAME="creative-assets",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",TEXT_OVERLAY_ENABLE_READINESS_PROBE="false",LOG_LEVEL="INFO"
priority=340
