# deployment/supervisord-allinone.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

# ============================================================================
# Supervisord Web UI - Access at http://localhost:9002
# ============================================================================
[inet_http_server]
port=*:9002
username=admin
password=admin

[supervisorctl]
serverurl=http://localhost:9002
username=admin
password=admin

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ============================================================================
# Infrastructure Services
# ============================================================================

[program:mongodb]
command=/usr/local/bin/mongod --dbpath=/data/mongodb --bind_ip_all --quiet --logpath=/dev/null
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/mongodb.log
stderr_logfile=/var/log/supervisor/mongodb-error.log
priority=10

[program:nats]
command=/usr/local/bin/nats-server -js --store_dir=/data/nats --addr=0.0.0.0 --http_port=8222
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/nats.log
stderr_logfile=/var/log/supervisor/nats-error.log
priority=20

[program:minio]
command=/usr/local/bin/minio server /data/minio --address=:9000 --console-address=:9001
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/minio.log
stderr_logfile=/var/log/supervisor/minio-error.log
environment=MINIO_ROOT_USER="minioadmin",MINIO_ROOT_PASSWORD="minioadmin"
priority=30

# ============================================================================
# Management UIs
# ============================================================================

[program:mongo-express]
command=/usr/local/bin/mongo-express
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/mongo-express.log
stderr_logfile=/var/log/supervisor/mongo-express-error.log
environment=ME_CONFIG_MONGODB_URL="mongodb://localhost:27017/",ME_CONFIG_MONGODB_ENABLE_ADMIN="true",ME_CONFIG_BASICAUTH_USERNAME="admin",ME_CONFIG_BASICAUTH_PASSWORD="admin",VCAP_APP_PORT="8081"
priority=100

# ============================================================================
# Application Services
# ============================================================================

[program:api]
command=python -m uvicorn src.api.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api-error.log
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",S3_ENDPOINT_URL="http://localhost:9000",S3_ACCESS_KEY_ID="minioadmin",S3_SECRET_ACCESS_KEY="minioadmin",S3_BUCKET_NAME="creative-assets",LOG_LEVEL="INFO"
priority=200

[program:web]
command=streamlit run src/web/app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/web.log
stderr_logfile=/var/log/supervisor/web-error.log
environment=API_BASE_URL="http://localhost:8000",LOG_LEVEL="INFO"
priority=210

[program:context-enricher]
command=python src/context_enricher/main.py
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/context-enricher.log
stderr_logfile=/var/log/supervisor/context-enricher-error.log
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",LOG_LEVEL="INFO"
priority=300

[program:creative-generator]
command=python src/creative_generator/main.py
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/creative-generator.log
stderr_logfile=/var/log/supervisor/creative-generator-error.log
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",LOG_LEVEL="INFO"
priority=310

[program:image-generator]
command=python src/image_generator/main.py
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/image-generator.log
stderr_logfile=/var/log/supervisor/image-generator-error.log
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",S3_ENDPOINT_URL="http://localhost:9000",S3_ACCESS_KEY_ID="minioadmin",S3_SECRET_ACCESS_KEY="minioadmin",S3_BUCKET_NAME="creative-assets",LOG_LEVEL="INFO"
priority=320

[program:brand-composer]
command=python src/brand_composer/main.py
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/brand-composer.log
stderr_logfile=/var/log/supervisor/brand-composer-error.log
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",S3_ENDPOINT_URL="http://localhost:9000",S3_EXTERNAL_ENDPOINT_URL="http://localhost:9000",S3_ACCESS_KEY_ID="minioadmin",S3_SECRET_ACCESS_KEY="minioadmin",S3_BUCKET_NAME="creative-assets",LOG_LEVEL="INFO"
priority=330

[program:text-overlay]
command=python src/text_overlay/main.py
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/text-overlay.log
stderr_logfile=/var/log/supervisor/text-overlay-error.log
environment=MONGODB_URL="mongodb://localhost:27017",MONGODB_DB_NAME="creative_campaign",NATS_URL="nats://localhost:4222",S3_ENDPOINT_URL="http://localhost:9000",S3_EXTERNAL_ENDPOINT_URL="http://localhost:9000",S3_ACCESS_KEY_ID="minioadmin",S3_SECRET_ACCESS_KEY="minioadmin",S3_BUCKET_NAME="creative-assets",LOG_LEVEL="INFO"
priority=340
